# resources/dummy.service
[Unit]
Description="dummy service to keep wsl awake"

[Service]
ExecStart="/usr/local/bin/dummy"
